# persoc.timer
[Unit]
Description=Schedule Persoc

[Timer]
OnUnitActiveSec=10s
AccuracySec=1s
Unit=persoc.service
Persistent=true

[Install]
WantedBy=timers.target
